# crontab
0 6 * * * /usr/local/bin/python /app/manage.py contact_daily
